# Discrete.service
[Unit]
Description=Discrete Monitoring Service
Wants=network-online.target
After=network.target network-online.target


[Service]
WorkingDirectory=/home/pi/Discrete_console_v2/bin
ExecStartPre=/bin/sleep 5
ExecStart=/home/pi/Discrete_console_v2/bin/Discrete_console
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=discrete
User=root
Group=root
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
